# systemd/qwen3-asr-ime.service
[Unit]
Description=Qwen3-ASR Voice Input Daemon
After=graphical-session.target qwen3-asr-server.service
Requires=qwen3-asr-server.service

[Service]
ExecStart=/home/aurora/miniconda3/bin/python3 -m qwen3_asr_ime.daemon.service
Restart=on-failure
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=default.target
