# systemd/qwen3-asr-server.service
[Unit]
Description=Qwen3-ASR Model Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{PYTHON}} {{PROJECT_DIR}}/tools/asr_server.py
Restart=on-failure
RestartSec=5
Environment="PYTHONUNBUFFERED=1"
Environment="QWEN3_ASR_MODE=offline"
Environment="QWEN3_ASR_BACKEND=transformers"
Environment="QWEN3_ASR_MODEL=/Data2/Models/Qwen3-ASR-1.7B"
Environment="QWEN3_ASR_DEVICE=auto"
Environment="QWEN3_ASR_PORT=8000"

[Install]
WantedBy=default.target
